# php/scheduled-subscription-charges.crontab
5 17 * * * /usr/local/bin/php /var/www/symfony_docker/bin/console app:subscription:charge > /tmp/subscription-charge.log 2> /tmp/subscription-charge.error.log
10 0 * * * /usr/local/bin/php /var/www/symfony_docker/bin/console app:subscription:charge > /tmp/subscription-charge.log 2> /tmp/subscription-charge.error.log
5 10 * * * /usr/local/bin/php /var/www/symfony_docker/bin/console app:subscription:remind > /tmp/subscription-remind.log 2> /tmp/subscription-remind.error.log
